# deployment/config_file/supervisord.conf
[program:yasbossapp]
command=/var/www/yesboss/deployment/gunicorn_start_pro
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/yesboss/logs/gunicorn-error.log